# systemd/uwb-mqtt-publisher.service
[Unit]
Description=UWB MQTT Publisher Service
Documentation=https://example.com/uwb-mqtt-publisher
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/uwb-mqtt-publisher
ExecStart=/usr/bin/uwb-mqtt-publisher ${UART_PORT} \
    --mqtt-broker ${MQTT_BROKER} \
    --mqtt-port ${MQTT_PORT} \
    --mqtt-topic ${MQTT_TOPIC} \
    --mqtt-rate-limit ${MQTT_RATE_LIMIT} \
    ${EXTRA_ARGS}
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=uwb-mqtt-publisher

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/dev

# Resource limits
LimitNOFILE=1024
MemoryLimit=64M

# Restart policy
StartLimitIntervalSec=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
